# files/etc/systemd/system/update-machines.timer
[Unit]
Description=Timer to poll /vagrant/shared/hosts

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
